# build/supervisord.conf
[supervisord]
idfile=/run/supervisord.pid
pidfile=/var/run/logsupervisord.pid
logfile=/dev/stdout
user=root

[program:app-swoole]
process_name=%(program_name)s
command=php -d variables_order=EGPCS /app/public/index.php
user=root
stopwaitsecs=30
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:app-job-worker]
process_name=%(program_name)s
command=php -d variables_order=EGPCS /app/bin/console.php jobs:process
user=root
stopwaitsecs=30
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:app-cron]
process_name=%(program_name)s
command=crond -f
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
